# deploy/bookstream.service
[Unit]
Description=Bookstream Next.js application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roomhacker
Group=roomhacker
WorkingDirectory=/home/roomhacker/PycharmProjects/bookstream
LimitNOFILE=65536
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
Environment=DATABASE_URL=file:/home/roomhacker/PycharmProjects/bookstream/db/custom.db
Environment=BOOKSTREAM_PUBLIC_DIR=/home/roomhacker/PycharmProjects/bookstream/public
Environment=BOOKSTREAM_MODE_FILE=/home/roomhacker/PycharmProjects/bookstream/.zscripts/service.mode
ExecStart=/home/roomhacker/PycharmProjects/bookstream/deploy/bookstream-run.sh
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
